# deploy/pdf_server_supervisor.ini
[program:pdf_server]
command=/usr/local/bin/uwsgi
  --socket /tmp/pdf_server.sock
  --home /var/www/reports/venv
  --pythonpath /var/www/reports/
  --wsgi-file /var/www/reports/src/wsgi.py
  --callable app
  --max-requests 1000
  --master
  --processes 4
  --chmod
directory=/var/www/reports
autostart=true
autorestart=true